# udev/rules.d/70-snap.color-picker.rules
# This file is automatically generated.
# opengl
KERNEL=="nvhost-*", TAG+="snap_color-picker_color-picker"
# opengl
KERNEL=="nvmap", TAG+="snap_color-picker_color-picker"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_color-picker_color-picker"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_color-picker_color-picker"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_color-picker_color-picker"
# opengl
KERNEL=="vchiq", TAG+="snap_color-picker_color-picker"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_color-picker_color-picker"
TAG=="snap_color-picker_color-picker", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_color-picker_color-picker $devpath $major:$minor"
