# udev/rules.d/70-snap.telegram-desktop.rules
# This file is automatically generated.
# opengl
KERNEL=="nvhost-*", TAG+="snap_telegram-desktop_telegram-desktop"
# opengl
KERNEL=="nvmap", TAG+="snap_telegram-desktop_telegram-desktop"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_telegram-desktop_telegram-desktop"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_telegram-desktop_telegram-desktop"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_telegram-desktop_telegram-desktop"
# opengl
KERNEL=="vchiq", TAG+="snap_telegram-desktop_telegram-desktop"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_telegram-desktop_telegram-desktop"
TAG=="snap_telegram-desktop_telegram-desktop", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_telegram-desktop_telegram-desktop $devpath $major:$minor"
