# systemd/system/snap-telegram\x2ddesktop-2198.mount
[Unit]
Description=Mount unit for telegram-desktop, revision 2198
Before=snapd.service

[Mount]
What=/var/lib/snapd/snaps/telegram-desktop_2198.snap
Where=/snap/telegram-desktop/2198
Type=squashfs
Options=nodev,ro,x-gdu.hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
